# uartserver.service
[Unit]
Description=UART serial communication 

[Service]
ExecStart=/usr/bin/uartserver

[Install]
WantedBy=multi-user.target
EOF